# integrations/systemd/crod-claude-hook.service
[Unit]
Description=CROD Claude Integration Hook
After=network.target

[Service]
Type=simple
User=daniel
WorkingDirectory=/home/daniel/Schreibtisch/Crod Programming/CROD-Helper-Member-7
ExecStart=/usr/bin/node /home/daniel/Schreibtisch/Crod Programming/CROD-Helper-Member-7/integrations/claude-imitation/crod-always-active.js
Restart=always
RestartSec=10
StandardOutput=append:/home/daniel/.claude/crod-learning.log
StandardError=append:/home/daniel/.claude/crod-learning-error.log

# Environment
Environment="NODE_ENV=production"
Environment="CROD_MODE=LEARNING"
Environment="CLAUDE_INTEGRATION=true"

[Install]
WantedBy=multi-user.target