# scripts/services/crod-always-active.service
[Unit]
Description=CROD Always Active Neural Network
After=network.target

[Service]
Type=simple
User=daniel
WorkingDirectory=/home/daniel/Schreibtisch/Crod Programming
ExecStart=/usr/bin/node /home/daniel/Schreibtisch/Crod Programming/crod-always-active.js
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal

# Environment
Environment="NODE_ENV=production"

[Install]
WantedBy=multi-user.target